# crontab
0 23 * * * find /tmp -type f -atime +1 -delete
